# scripts/upstart/truebot.conf
# upstart config

description "TrueBot 2.0"
author "Main Author: Trooble - Scripts: Slurpee"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
exec /home/ubuntu/truebot
